# DEBIAN/postinst
#!/bin/sh
update-rc.d -f bsprinter defaults
/etc/init.d/bsprinter restart