# DEBIAN/postinst
#!/bin/sh
update-rc.d -f bs_printer defaults
/etc/init.d/bs_printer restart